# codedeploy/ec2/entrypoints/supervisord.queue.conf
[supervisord]
nodaemon=true

[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/php /var/www/html/artisan queue:work database --sleep=3 --tries=3
autostart=true
autorestart=true

# 컨테이너 하나에서 5개의 프로세스를 동작하게 하는데... 검토좀 해보자
# numprocs=5 
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0